# debs/AppSync14/data/var/lib/dpkg/info/net.angelxwind.appsyncunified.list
/.
/usr
/usr/bin
/usr/bin/asu_inject
/Library
/Library/MobileSubstrate
/Library/MobileSubstrate/DynamicLibraries
/Library/MobileSubstrate/DynamicLibraries/AppSyncUnified.dylib
/Library/MobileSubstrate/DynamicLibraries/AppSyncUnified.plist
/Library/LaunchDaemons
/Library/LaunchDaemons/net.angelxwind.asu_inject.plist
